# dot_config/systemd/user/battery-check.service
[Unit]
Description=Check battery level and notify

[Service]
Type=oneshot
ExecStart=/home/ashutosh/.config/scripts/battery-check
