# dot_config/systemd/user/waybar-ycal.service
[Unit]
Description=Waybar Google Calendar popup daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 $INSTALL_DIR/popup.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
